# debian-rfs-stats.cron
0 1 * * * daves (cd /home/daves/github/debian-rfs-stats; ./bin/daily >/dev/null)
